# eventhandler/dist/debian/DEBIAN/postrm
#!/bin/bash

[ -e /etc/init.d/etoad-manager ] && update-rc.d etoad-manager remove
